# crontab
* * * * * bash /pathtofile/auto-start servicename > /dev/null
